# atop-rotate.timer
[Unit]
Description=Daily atop restart

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
